# infra/systemd/gemma_bot.service
[Unit]
Description=Gemma Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/gemma_bot
Environment=GEMMA_PROJECT_ROOT=/srv/gemma_bot
EnvironmentFile=-/srv/gemma_bot/.env
ExecStart=/srv/gemma_bot/venv/bin/python3 /srv/gemma_bot/main.py
Restart=on-failure
RestartSec=10
# Логи: panel_nohup_bot.log при direct; journalctl -u gemma_bot
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
